# AI学习/02-工具使用/launchd-configs/com.user.ai-news-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.ai-news-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>sleep 20; cd /Users/zhqznc/Documents/新闻 &amp;&amp; /usr/local/bin/claude --dangerously-skip-permissions -p "收集今天的AI新闻总结成文档"</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zhqznc/Documents/新闻/logs/ai-news-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zhqznc/Documents/新闻/logs/ai-news-error.log</string>

    <key>AbandonProcessGroup</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>
</dict>
</plist>